# DEBIAN/postinst
#!/bin/sh
set -e

# Debconf-Modul laden
. /usr/share/debconf/confmodule

CONF="/etc/vdr/conf.d/vdr-rectools.conf"

# Verzeichnis sicherstellen
mkdir -p /etc/vdr/conf.d

if [ ! -f "$CONF" ]; then
    cat <<EOFF > "$CONF"
MAX_FILES=10
MIN_FREE_GB=20
AUTO_CLEAN=1
IMPORT_DIR="/srv/videos-other"
PAUSE_WORK=30
PAUSE_CHECK=2
SNAPSHOT_TIME="00:05:00"
EOFF
fi

db_get vdr-rectools/mail
MAIL_NOTIFY="$RET"

db_get vdr-rectools/timer
if [ "$RET" = "true" ]; then
    AUTO_TIMER=1
else
    AUTO_TIMER=0
fi

db_get vdr-rectools/subtitles
if [ "$RET" = "true" ]; then
    AUTO_SUB_DOWNLOAD=1
else
    AUTO_SUB_DOWNLOAD=0
fi

sed -i "/^MAIL_NOTIFY=/d" "$CONF"
echo "MAIL_NOTIFY=\"$MAIL_NOTIFY\"" >> "$CONF"

sed -i "/^AUTO_TIMER=/d" "$CONF"
echo "AUTO_TIMER=$AUTO_TIMER" >> "$CONF"

sed -i "/^AUTO_SUB_DOWNLOAD=/d" "$CONF"
echo "AUTO_SUB_DOWNLOAD=$AUTO_SUB_DOWNLOAD" >> "$CONF"

sed -i "/^SUB_LANG=/d" "$CONF"
echo "SUB_LANG=\"de\"" >> "$CONF"

if [ -f "/var/lib/vdr/reccmds.conf" ] || [ -d "/usr/share/yavdr-ansible" ]; then
    mkdir -p /etc/systemd/system/vdr.service.d
    cat <<EOFF > /etc/systemd/system/vdr.service.d/99-rectools-menu.conf
[Service]
ExecStartPre=+/bin/bash -c "grep -q repair_single /var/lib/vdr/reccmds.conf || echo \"Aufnahme reparieren (Rectools) : /usr/bin/vdr-rectools repair_single\" >> /var/lib/vdr/reccmds.conf"
ExecStartPre=+/bin/bash -c "grep -q cut_single /var/lib/vdr/reccmds.conf || echo \"Werbung schneiden (Rectools)   : /usr/bin/vdr-rectools cut_single\" >> /var/lib/vdr/reccmds.conf"
ExecStartPre=+/bin/bash -c "grep -q shrink_single /var/lib/vdr/reccmds.conf || echo \"Platz sparen H.265 (Rectools)  : /usr/bin/vdr-rectools shrink_single\" >> /var/lib/vdr/reccmds.conf"
ExecStartPre=+/bin/bash -c "grep -q sync_single /var/lib/vdr/reccmds.conf || echo \"Plex/Kodi Sync (Rectools)      : /usr/bin/vdr-rectools sync_single\" >> /var/lib/vdr/reccmds.conf"
EOFF
else
    rm -f /etc/systemd/system/vdr.service.d/99-rectools-menu.conf
fi

if [ -d /run/systemd/system ]; then
    systemctl daemon-reload
    systemctl enable --now vdr-rectools.timer || true
fi

db_stop

#DEBHELPER#

exit 0
