# debian/postinst
#!/bin/sh
set -e

# Debconf-Modul laden
. /usr/share/debconf/confmodule

CONF="/etc/vdr/conf.d/vdr-rectools.conf"
mkdir -p /etc/vdr/conf.d

# 1. Grundgerüst erstellen, falls Datei fehlt
if [ ! -f "$CONF" ]; then
    cat <<EOFF > "$CONF"
MAX_FILES=10
MIN_FREE_GB=20
AUTO_CLEAN=1
IMPORT_DIR="/srv/video/Filme"
PAUSE_WORK=30
PAUSE_CHECK=2
SNAPSHOT_TIME="00:05:00"
MAIL_NOTIFY=""
AUTO_TIMER=0
AUTO_START_NIGHT=0
AUTO_SUB_DOWNLOAD=1
SUB_LANG="de"
USE_TVSCRAPER=0
TVSCRAPER_MODE="immediate"
TVSCRAPER_DB="/var/cache/vdr/plugins/tvscraper/tvscraper.sqlite"
EOFF
else
    # 1b. Bestehende Datei prüfen und fehlende Variablen nachrüsten
    for VAR in MAX_FILES MIN_FREE_GB AUTO_CLEAN IMPORT_DIR PAUSE_WORK PAUSE_CHECK SNAPSHOT_TIME MAIL_NOTIFY AUTO_TIMER AUTO_START_NIGHT AUTO_SUB_DOWNLOAD SUB_LANG USE_TVSCRAPER TVSCRAPER_MODE TVSCRAPER_DB; do        
        if ! grep -q "^${VAR}=" "$CONF"; then
            case $VAR in
                AUTO_START_NIGHT)  echo "AUTO_START_NIGHT=0" >> "$CONF" ;;
                AUTO_TIMER)        echo "AUTO_TIMER=0" >> "$CONF" ;;
                AUTO_SUB_DOWNLOAD) echo "AUTO_SUB_DOWNLOAD=1" >> "$CONF" ;;
                SUB_LANG)          echo "SUB_LANG=\"de\"" >> "$CONF" ;;
                MAIL_NOTIFY)       echo "MAIL_NOTIFY=\"\"" >> "$CONF" ;;
                CRF_H264_DEFAULT)  echo "CRF_H264_DEFAULT=23" >> "$CONF" ;;
                PRESET_H264_DEFAULT) echo "PRESET_H264_DEFAULT=\"medium\"" >> "$CONF" ;;
                CRF_H265_DEFAULT)  echo "CRF_H265_DEFAULT=23" >> "$CONF" ;;
                PRESET_H265_DEFAULT) echo "PRESET_H265_DEFAULT=\"medium\"" >> "$CONF" ;;
                CRF_H264_FALLBACK) echo "CRF_H264_FALLBACK=23" >> "$CONF" ;;
                PRESET_H264_FALLBACK) echo "PRESET_H264_FALLBACK=\"fast\"" >> "$CONF" ;;
                MIN_COMPRESSION_RATIO_H264) echo "MIN_COMPRESSION_RATIO_H264=70" >> "$CONF" ;;
                MIN_COMPRESSION_RATIO_H265) echo "MIN_COMPRESSION_RATIO_H265=50" >> "$CONF" ;;
                MIN_COMPRESSION_RATIO_H264_FALLBACK) echo "MIN_COMPRESSION_RATIO_H264_FALLBACK=70" >> "$CONF" ;;
                HW_ACCEL)          echo "HW_ACCEL=\"none\"" >> "$CONF" ;;
                AUTO_ENCODE_IMPORT) echo "AUTO_ENCODE_IMPORT=1" >> "$CONF" ;;
                USE_TVSCRAPER)     echo "USE_TVSCRAPER=0" >> "$CONF" ;;
                TVSCRAPER_MODE)    echo "TVSCRAPER_MODE=\"immediate\"" >> "$CONF" ;;
                TVSCRAPER_DB)      echo "TVSCRAPER_DB=\"/var/cache/vdr/plugins/tvscraper/tvscraper.sqlite\"" >> "$CONF" ;;
                *) echo "$VAR=\"\"" >> "$CONF" ;;
            esac
        fi
    done
fi

# 2. Debconf-Werte auslesen
db_get vdr-rectools/mail
MAIL_VAL="$RET"
db_get vdr-rectools/timer
[ "$RET" = "true" ] && TIMER_VAL=1 || TIMER_VAL=0
db_get vdr-rectools/subtitles
[ "$RET" = "true" ] && SUB_VAL=1 || SUB_VAL=0

# Werte sicher in der Datei aktualisieren
sed -i "/^MAIL_NOTIFY=/c\MAIL_NOTIFY=\"$MAIL_VAL\"" "$CONF"
sed -i "/^AUTO_TIMER=/c\AUTO_TIMER=$TIMER_VAL" "$CONF"
sed -i "/^AUTO_SUB_DOWNLOAD=/c\AUTO_SUB_DOWNLOAD=$SUB_VAL" "$CONF"

# 3. Systemd-Teil
if [ -d /run/systemd/system ]; then
    systemctl daemon-reload
    systemctl enable vdr-rectools.timer || true
    systemctl start vdr-rectools.timer || true
fi

# 4. VDR OSD-Befehle (drop-in)
if [ -f "/var/lib/vdr/reccmds.conf" ] || [ -d "/usr/share/yavdr-ansible" ]; then
    mkdir -p /etc/systemd/system/vdr.service.d
    cat <<EOFF > /etc/systemd/system/vdr.service.d/99-rectools-menu.conf
[Service]
ExecStartPre=+/bin/bash -c "grep -q repair_single /var/lib/vdr/reccmds.conf || echo \"Aufnahme reparieren (Rectools) : /usr/bin/vdr-rectools repair_single\" >> /var/lib/vdr/reccmds.conf"
ExecStartPre=+/bin/bash -c "grep -q cut_single /var/lib/vdr/reccmds.conf || echo \"Werbung schneiden (Rectools)   : /usr/bin/vdr-rectools cut_single\" >> /var/lib/vdr/reccmds.conf"
ExecStartPre=+/bin/bash -c "grep -q shrink_single /var/lib/vdr/reccmds.conf || echo \"Platz sparen H.265 (Rectools)  : /usr/bin/vdr-rectools shrink_single\" >> /var/lib/vdr/reccmds.conf"
ExecStartPre=+/bin/bash -c "grep -q sync_single /var/lib/vdr/reccmds.conf || echo \"Plex/Kodi Sync (Rectools)      : /usr/bin/vdr-rectools sync_single\" >> /var/lib/vdr/reccmds.conf"
EOFF
fi

# 5. RECHTE UND ARBEITSVERZEICHNISSE (FIX FÜR USER VDR)
# Erstellt staging und repaired_files, korrigiert Rechte für video.00
mkdir -p /srv/vdr/tmp/staging /srv/vdr/repaired_files /srv/vdr/video.00
chown -R vdr:vdr /srv/vdr/tmp /srv/vdr/repaired_files /srv/vdr/video.00
chmod -R 775 /srv/vdr/tmp /srv/vdr/repaired_files /srv/vdr/video.00

# Logfile für vdr-user schreibbar machen
touch /var/log/vdr-rectools.log
chown vdr:vdr /var/log/vdr-rectools.log
chmod 664 /var/log/vdr-rectools.log

db_stop
#DEBHELPER#
exit 0
